# src/systemd/collecty.service.in
[Unit]
Description=collecty - A system data collecting daemon

[Service]
ExecStart=@bindir@/collectyd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
