# src/systemd/collecty.service.in
[Unit]
Description=collecty - A system data collecting daemon

[Service]
Type=dbus
ExecStart=@sbindir@/collectyd
ExecReload=/bin/kill -HUP $MAINPID
BusName=org.ipfire.collecty1

[Install]
WantedBy=multi-user.target
